# deploy/lina.service
[Unit]
Description=Lina web application
After=network.target

[Service]
Type=simple
User=lina
Group=lina
WorkingDirectory=/opt/lina/current
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3000
Environment=LINA_DATA_DIR=/opt/lina/data
EnvironmentFile=-/etc/lina.env
ExecStart=/opt/node/bin/node /opt/lina/current/server.js
Restart=always
RestartSec=2
TimeoutStopSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/lina

[Install]
WantedBy=multi-user.target
