# deploy/lina-telegram-bot.service
[Unit]
Description=Lina Telegram bot
After=network-online.target lina.service
Wants=network-online.target lina.service

[Service]
Type=simple
User=lina-bot
Group=lina-bot
EnvironmentFile=/etc/lina-telegram-bot.env
ExecStart=/usr/bin/python3 /opt/lina/current/lina_telegram_bot.py
Restart=always
RestartSec=3
TimeoutStopSec=40
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
UMask=0027
TasksMax=64

[Install]
WantedBy=multi-user.target
