# deploy/lina.service
[Unit]
Description=Lina web application
After=network.target

[Service]
Type=simple
User=lina
Group=lina
WorkingDirectory=/opt/lina/current
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3000
EnvironmentFile=-/etc/lina.env
ExecStart=/opt/node/bin/node /opt/lina/current/server.js
Restart=always
RestartSec=2
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
UMask=0027
TasksMax=256
LimitNOFILE=65536
ReadWritePaths=/opt/lina/current/.next/cache

[Install]
WantedBy=multi-user.target
